# debian/hdparm.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="hd[a-z]", \
	RUN+="/lib/udev/hdparm"
