# debian/hdparm.postinst
#!/bin/sh

# Clean up cdbs damage
rm -f /etc/udev/rules.d/50-hdparm.rules


if dpkg --compare-versions "$2" lt "8.6-1ubuntu1"; then
	update-rc.d -f hdparm remove
fi

if [ -n "$2" ] && dpkg --compare-versions "$2" lt 8.9-2; then
  rm -fv /etc/rc0.d/K75hdparm /etc/rc6.d/K75hdparm
fi
#DEBHELPER#
