# debian/hdparm.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="hd[a-z]", \
	RUN+="/etc/init.d/hdparm hotplug"
